# templates/sysbench_prepared.target
[Unit]
Description=Sysbench has ran prepare command
Requires=multi-user.target network.target
Conflicts=rescue.service rescue.target
AllowIsolate=yes